# bin/notify-send.cron
#!/bin/bash
# This script is in public domain.
# Author:  Leonid Evdokimov <leon@darkk.net.ru>

env.session notify-send "$@"
