# .config/upstart/redshift.conf
# Part of ~/bin/update-geolocation and ~/bin/geolocate glue.
#
# -- Leonid Evdokimov <leon@darkk.net.ru>

start on started gnome-session

respawn

exec redshift -l $(cat ~/.cache/geolocate/last)
